# roles/interesting-repositories/templates/gickup-sync.service.j2
[Unit]
Description=Synchronize Gickup repositories

[Service]
ExecStart=/usr/bin/gickup-sync 
User=wynand
Group=wynand
